# deploy/expense-tracker.service
[Unit]
Description=Expense Tracker (backend API + built dashboard, single process)
After=network.target

[Service]
WorkingDirectory=%h/expense_tracker/backend
ExecStart=%h/expense_tracker/backend/.venv/bin/uvicorn app.presentation.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
